# files/init/tkellem.conf
description 'tkellem'

start on runlevel [2345]
stop on runlevel [016]

console log

setuid tkellem
setgid tkellem

env HOME=/var/lib/tkellem

respawn
exec /usr/local/bin/tkellem run

